# templates/postgresql/pg_backup.timer.erb
[Unit]
Description=Daily backup of PostgreSQL database

[Timer]
OnCalendar=*-*-* <%= @backup_time %>
Persistent=true

[Install]
WantedBy=timers.target
