# templates/postgresql/pg_backup.service.erb
[Unit]
Description=Daily backup of PostgreSQL database
Requires=postgresql.service

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart=<%= @backup_dir %>/pg_backup_rotated.sh

[Install]
WantedBy=network.target
