# templates/gitea/gitea-backup.service.erb
[Unit]
Description=Backup gitea
Requires=docker-gitea.service

[Service]
Type=oneshot
ExecStartPre=find <%= @gitea_backup_dir %> -type f -mtime +30 -exec rm -f {} \;
ExecStart=/usr/bin/docker exec -i -w /data/backup  gitea bash -c '/app/gitea/gitea dump -c /data/custom/conf/app.ini'

[Install]
WantedBy=network.target
