# templates/gitea/gitea-backup.timer.erb
[Unit]
Description=Backup gitea

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
